# code/code.service
[Unit]
Description=Gunicorn instance to serve code
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/app
Environment="PATH=/home/ubuntu/app/myenv/bin"
ExecStart=/home/ubuntu/app/myenv/bin/gunicorn --workers 3 --bind 0.0.0.0:5000 -m 007 wsgi:app

[Install]
WantedBy=multi-user.target